# z-backend/trudy-backend.service
[Unit]
Description=Trudy Backend API Service
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/backend
Environment="PATH=/opt/backend/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/opt/backend/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 4
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=trudy-backend

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/backend /mnt/storage

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
